# .planning/phases/13-deployment/freediving.service.example
[Unit]
Description=Fishly — Freediving Breathhold Trainer
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/var/www/freediving
Environment=PATH=/root/.nvm/versions/node/v20.20.1/bin:/usr/local/bin:/usr/bin:/bin
EnvironmentFile=/var/www/freediving/.env.production
ExecStart=/var/www/freediving/start_freediving.sh
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
